# recipes-devtools/cannelloni-server/files/spawn-cannelloni-server.service
[Unit]
Description=Cannelloni Client Service
After=network.target

[Service]
Type=simple
ExecStart=/opt/uthp/scripts/cannelloni-server/start.sh
User=root

[Install]
WantedBy=multi-user.target